# rc.local
#!/bin/sh -e
#
# rc.local
# this needs to end with an exit 0

/usr/bin/sudo -u pi /home/pi/startup.sh > /home/pi/log 2>&1 &
exit 0
